# root/etc/services.d/dropbox-monitor/run
#!/usr/bin/with-contenv bash

export HOME=/dropbox
export LD_PRELOAD=/usr/local/lib/libdropbox_fs_fix.so
s6-setuidgid abc /usr/local/bin/dropbox status
sleep 30
