# root/etc/services.d/dropbox/run
#!/usr/bin/with-contenv bash

export HOME=/dropbox
export LD_PRELOAD=/usr/local/lib/libdropbox_fs_fix.so
exec \
	s6-setuidgid abc /dropbox/.dropbox-dist/dropboxd "${RUN_OPTS}"
